# packages/pscripts/pscripts-0.1.100.tar.gz/pscripts-0.1.100/config/update_external_ip.service
[Unit]
Description=Update your dynamic IP address with a Dynamic DNS service

[Service]
ExecStart=/usr/bin/update_external_ip

[Install]
WantedBy=multi-user.target